# data/powertop.service
[Unit]
Description=Powertop startup
Before=powertop_tunables.service

[Service]
ExecStart=/usr/sbin/powertop --auto-tune
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
